# daemon/opencode-user@.service
[Unit]
Description=OpenCode per-user daemon for %i
After=network.target

[Service]
Type=simple
User=%i
ExecStart=/usr/local/bin/opencode serve --unix-socket /run/user/%U/opencode/daemon.sock
Restart=on-failure
RestartSec=3
RuntimeDirectory=opencode
RuntimeDirectoryMode=0700
Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=default.target
