# daemon/opencode-gateway.service
[Unit]
Description=TheSmartAI Gateway (C root daemon)
Documentation=https://github.com/sst/opencode
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/opencode-gateway
Restart=on-failure
RestartSec=3
KillMode=mixed
TimeoutStopSec=10

# All config from single source of truth
EnvironmentFile=/etc/opencode/opencode.cfg

# Runtime directory for ephemeral files (PID tracking, etc.)
RuntimeDirectory=opencode-gateway
RuntimeDirectoryMode=0755

# Persistent state directory for JWT secret (survives service restarts)
StateDirectory=opencode-gateway
StateDirectoryMode=0700

# Security hardening (gateway needs root for PAM + setuid, but limit the rest)
NoNewPrivileges=no
ProtectSystem=strict
ReadWritePaths=/run/opencode-gateway /var/lib/opencode-gateway /run/user /tmp /etc/opencode /opt/opencode-apps /var/log/opencode
ProtectHome=no

[Install]
WantedBy=multi-user.target
